# src/main/docker/liberty.conf
[program:liberty]
command=/opt/ibm/wlp/bin/server start
